# stow/autostart/.config/autostart/autostart.desktop
[Desktop Entry]
Name=Autostart
Exec=sh /home/f1p/git/dotfiles/stow/autostart/bin/autostart.sh
Type=Application
